# roles/ec2-ami-base/files/etc.systemd.docker.brtfs.mounts.service
[Unit]
Description=Mounts BRTFS docker mountpoints

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-btrfs-mounts.sh
RemainAfterExit=true

[Install]
WantedBy=docker.service
