# roles/ec2-ami-base/files/etc.init.docker.brtfs.mounts.conf
# one time script

start on starting docker

post-start script
	/usr/local/bin/docker-btrfs-mounts.sh
end script
